# etc/systemd/system/certbot-renew-all.service
[Unit]
Description=Renew all certs generated by letsencrypts certbot

[Service]
Type=oneshot
ExecStartPre=/bin/systemctl stop nginx
ExecStart=/home/pi/dev/certbot/certbot-auto renew --quiet
ExecStartPost=/bin/systemctl start nginx
