# etc/systemd/system/mongod.service
[Unit]
Description=MongoDB Daemon

[Service]
ExecStart=/usr/bin/mongod
ExecStop=kilall /usr/bin/mongod

[Install]
WantedBy=multi-user.target
